# files/etc/uci-defaults/99-custom
# Beware! This script will be in /rom/etc/uci-defaults/ as part of the image.
# Uncomment lines to apply:
#
# wlan_name="OpenWrt"
# wlan_password="12345678"
#
# root_password=""
lan_ip_address="192.168.0.1"
#
# pppoe_username=""
# pppoe_password=""

# log potential errors
exec >/tmp/setup.log 2>&1

if [ -n "$root_password" ]; then
  (echo "$root_password"; sleep 1; echo "$root_password") | passwd > /dev/null
fi

# Configure LAN interface
# More options: https://openwrt.org/docs/guide-user/base-system/basic-networking
if [ -n "$lan_ip_address" ]; then
  uci set network.lan.ipaddr="$lan_ip_address"
  uci commit network
fi
uci set network.lan.hostname='OpenWrt'
uci set network.lan.delegate='0'
uci set network.lan.ip6assign='64'
uci add_list network.lan.ip6class='local'
uci add_list network.lan.ip6class='wan6'
uci commit network

# Configure WAN interface
uci set network.wan.ipv6='1'
uci commit network

# Configure WAN6 interface
uci set network.wan6.device='@wan'
uci set network.wan6.reqaddress='try'
uci set network.wan6.reqprefix='auto'
uci commit network

# Configure WLAN
# More options: https://openwrt.org/docs/guide-user/network/wifi/basic#wi-fi_interfaces
if [ -n "$wlan_name" -a -n "$wlan_password" -a ${#wlan_password} -ge 8 ]; then
  uci set wireless.@wifi-device[0].country='MY'
  uci set wireless.@wifi-device[0].disabled='0'
  uci set wireless.@wifi-iface[0].disabled='0'
  uci set wireless.@wifi-iface[0].encryption='psk2+ccmp'
  uci set wireless.@wifi-iface[0].ssid="${wlan_name}-2.4G@unifi"
  uci set wireless.@wifi-iface[0].key="$wlan_password"
  # 5 GHz
  uci set wireless.@wifi-device[1].country='MY'
  uci set wireless.@wifi-device[1].disabled='0'
  uci set wireless.@wifi-iface[1].disabled='0'
  uci set wireless.@wifi-iface[1].encryption='psk2+ccmp'
  uci set wireless.@wifi-iface[1].ssid="${wlan_name}-5G@unifi"
  uci set wireless.@wifi-iface[1].key="$wlan_password"
  uci set wireless.@wifi-iface[1].wds='1'
  uci commit wireless
fi

# Configure PPPoE
# More options: https://openwrt.org/docs/guide-user/network/wan/wan_interface_protocols#protocol_pppoe_ppp_over_ethernet
if [ -n "$pppoe_username" -a "$pppoe_password" ]; then
  uci set network.wan.proto=pppoe
  uci set network.wan.username="$pppoe_username"
  uci set network.wan.password="$pppoe_password"
  uci set network.wan.device="wan.500"
  uci commit network
fi

# Configure Timezone & hostname
timezone='<+08>-8'
zonename='Asia/Kuala Lumpur'
uci set system.@system[0].timezone="$timezone"
uci set system.@system[0].zonename="$zonename"
uci delete system.ntp.server
uci add_list system.ntp.server='my.pool.ntp.org'
uci add_list system.ntp.server='sg.pool.ntp.org'
uci set system.@system[0].hostname='OpenWrt'
uci set network.lan.hostname="`uci get system.@system[0].hostname`"
uci commit

# Configure DNS provider
# Configure dnsmasq
uci -q delete dhcp.lan.dhcp_option
uci add_list dhcp.lan.dhcp_option="6,8.8.8.8,8.8.4.4"
uci commit dhcp
service dnsmasq restart
 
# Configure odhcpd
uci -q delete dhcp.lan.dns
uci add_list dhcp.lan.dns="2001:4860:4860::8888"
uci add_list dhcp.lan.dns="2001:4860:4860::8844"
uci commit dhcp
service odhcpd restart

echo "All done!"
