# files/etc/uci-defaults/99-custom
# Beware! This script will be in /rom/etc/uci-defaults/ as part of the image.
# Uncomment lines to apply:
#
# wlan_name="OpenWrt"
# wlan_password="12345678"
#
# root_password=""
lan_ip_address="192.168.0.2"
#
# pppoe_username=""
# pppoe_password=""

# log potential errors
exec >/tmp/setup.log 2>&1

if [ -n "$root_password" ]; then
  (echo "$root_password"; sleep 1; echo "$root_password") | passwd > /dev/null
fi

# Configure LAN
# More options: https://openwrt.org/docs/guide-user/base-system/basic-networking
if [ -n "$lan_ip_address" ]; then
  uci set network.lan.ipaddr="$lan_ip_address"
  uci commit network
fi

# Configure WLAN
# More options: https://openwrt.org/docs/guide-user/network/wifi/basic#wi-fi_interfaces
if [ -n "$wlan_name" -a -n "$wlan_password" -a ${#wlan_password} -ge 8 ]; then
  uci set wireless.@wifi-device[0].disabled='0'
  uci set wireless.@wifi-iface[0].disabled='0'
  uci set wireless.@wifi-iface[0].encryption='psk2'
  uci set wireless.@wifi-iface[0].ssid="${wlan_name}-5G@unifi"
  uci set wireless.@wifi-iface[0].key="$wlan_password"
  uci set wireless.@wifi-device[0].country='MY'
  uci set wireless.@wifi-iface[0].mode='sta'
  uci set wireless.@wifi-device[0].channel="36"
  uci set wireless.@wifi-device[0].htmode="VHT40"
  uci set wireless.@wifi-iface[0].wds='1'
  uci commit wireless
fi

# Configure PPPoE
# More options: https://openwrt.org/docs/guide-user/network/wan/wan_interface_protocols#protocol_pppoe_ppp_over_ethernet
if [ -n "$pppoe_username" -a "$pppoe_password" ]; then
  uci set network.wan.proto=pppoe
  uci set network.wan.username="$pppoe_username"
  uci set network.wan.password="$pppoe_password"
  uci commit network
fi

# Configure Timezone & hostname
timezone='<+08>-8'
zonename='Asia/Kuala Lumpur'
uci set system.@system[0].timezone="$timezone"
uci set system.@system[0].zonename="$zonename"
uci delete system.ntp.server
uci add_list system.ntp.server='my.pool.ntp.org'
uci add_list system.ntp.server='sg.pool.ntp.org'
uci set system.@system[0].hostname='DIR842'
uci set network.lan.hostname="`uci get system.@system[0].hostname`"
uci commit

# Configure additional config for wds
uci set network.lan.netmask='255.255.255.0'
uci set network.lan.gateway='192.168.0.1'
uci set network.lan.dns='192.168.0.1'
uci commmit

# Bridge WAN port with br-lan
#uci add_list network.@device[0].ports='eth0.2'
uci set network.@switch_vlan[0].ports='0t 4 3 2 1 5'
uci set network.@switch_vlan[1].ports='0t'
uci commit

# Delete wan & wan6 interface
uci del network.wan
uci del network.wan6
uci commit

# Disable IPv6
uci set network.@device[0].ipv6='0'
uci commit

echo "All done!"
