# files/etc/uci-defaults/99-custom
# Beware! This script will be in /rom/etc/uci-defaults/ as part of the image.
# Uncomment lines to apply:
#
# wlan_name="OpenWrt"
# wlan_password="12345678"
#
# root_password=""
lan_ip_address="192.168.0.2"
#
# pppoe_username=""
# pppoe_password=""

# log potential errors
exec >/tmp/setup.log 2>&1

if [ -n "$root_password" ]; then
  (echo "$root_password"; sleep 1; echo "$root_password") | passwd > /dev/null
fi

# network
uci set network.lan.ipaddr="$lan_ip_address"
uci set network.lan.gateway='192.168.0.1'
uci set network.lan.ip6assign='64'
uci del network.lan.dns
uci add_list network.lan.dns='192.168.0.1'
uci del network.@device[0].ports
uci set network.@device[0].ipv6='0'
uci add_list network.@device[0].ports='eth0.1'
uci add_list network.@device[0].ports='eth0.2'
uci commit network

# wireless
uci set  wireless.@wifi-device[0].country='MY'
uci set wireless.@wifi-device[0].disabled='0'

uci set wireless.@wifi-iface[0].mode='sta'
uci set wireless.@wifi-iface[0].ssid="${wlan_name}-5G@unifi"
uci set wireless.@wifi-iface[0].encryption='psk2+aes'
uci set wireless.@wifi-iface[0].key="$wlan_password"
uci set wireless.@wifi-iface[0].wds='1'
uci commit wireless

# dhcp
uci set dhcp.lan.ignore='1'
uci del dhcp.lan.ra
uci del dhcp.lan.dhcpv6
uci del dhcp.lan.ra_flags
uci commit dhcp

# system

uci set system.@system[0].hostname='DIR842-STA'
uci set system.@system[0].zonename='Asia/Kuala Lumpur'
uci set system.@system[0].timezone='<+08>-8'
uci del system.ntp.server
uci add_list system.ntp.server='my.pool.ntp.org'
uci add_list system.ntp.server='sg.pool.ntp.org'
uci commit system

echo "All done!"
